# contrib/oled-notifier.service
[Unit]
Description=OLED Display Notifier

[Service]
ExecStart=/usr/bin/oled-notifier

[Install]
WantedBy=multi-user.target
